# deploy/systemd/cftv-dashboard.service
[Unit]
Description=CFTV Dashboard Web Server
After=network.target

[Service]
Type=simple
User=vboxuser
WorkingDirectory=/home/vboxuser/SIMULA-O-CFTV
ExecStart=/home/vboxuser/SIMULA-O-CFTV/.venv/bin/python servidor_web_cftv.py --host 0.0.0.0 --port 8080
Restart=always
RestartSec=3

# Hardening defaults for lab/VM usage
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/home/vboxuser/SIMULA-O-CFTV

[Install]
WantedBy=multi-user.target
